# cmd/gitwww/gitwww.service
# /etc/systemd/system/gitwww.service
[Unit]
Description=GitWWW: automatic static-site builder
Documentation=exec:/var/opt/gitwww/gitwww -d
After=network.target docker.service
Requires=docker.service

# Don't use Environment= or EnvironmentFile= for credentials / secrets.
# But use LoadCredential=, LoadCredentialEncrypted= or SetCredentialEncrypted=
# Note: env vars set for a unit are exposed to unprivileged clients via D-Bus IPC.
# Moreover, env vars are propagated down the process tree,
# including across security boundaries (such as setuid/setgid executables),
# and hence might leak to processes that should not have access to the secret data.
# https://serverfault.com/a/1110351

[Service]
Type=simple
User=root
Environment="GITWWW_CFG=/var/opt/gitwww"
Environment="GITWWW_WWW=www"
Environment="GITWWW_LOG=INFO"
ExecStart=/var/opt/gitwww/gitwww

# -------------------------------------------------
# CPU scheduling – balanced trade‑off
# -------------------------------------------------
# Run with the idle scheduling policy so the daemon only gets CPU
# time when no higher‑priority task is runnable. When the system is
# otherwise idle, the process can use all available CPUs, giving
# BuildKit the full compute power it needs.
CPUSchedulingPolicy=idle

# Keep the traditional nice value at the lowest priority as an extra
# safeguard for systems that do not fully honor SCHED_IDLE.
Nice=19

# -------------------------------------------------
# Restart policy
# -------------------------------------------------
# Restart only on unexpected failures; a manual “systemctl stop” will keep it stopped.
Restart=on-failure
RestartSec=10

[Install]
# When the service is enabled (systemctl enable gitwww.service),
# systemd creates a symlink from this unit into
# /etc/systemd/system/multi-user.target.wants/.
# multi-user.target represents the normal, non‑graphical multi‑user
# system state (runlevel 3).  At boot, when the system reaches that
# target, all units wanted by it—including this one—are started
# automatically.
WantedBy=multi-user.target
